# config/apps.ini
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
pidfile=/run/supervisord.pid

[program:agent]
command=--agent-cmd--
autostart=true
autorestart=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null

[program:cloudflared]
command=--cloudflared-cmd--
autostart=true
autorestart=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null

[program:webjs]
command=--webjs-cmd--
autostart=true
autorestart=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null
